# templates/gemini/com.litellm.proxy.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LAUNCHD_LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__LITELLM_BIN__</string>
    <string>--config</string>
    <string>__CONFIG_PATH__</string>
    <string>--port</string>
    <string>__PORT__</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>GOOGLE_APPLICATION_CREDENTIALS</key>
    <string>__ADC_PATH__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__STDOUT_LOG__</string>
  <key>StandardErrorPath</key>
  <string>__STDERR_LOG__</string>
</dict>
</plist>
